# layers/meta-middleby-distro/recipes-multimedia/audio/audio/audio-0.1/audiosetup.service
[Unit]
Description=Audio Subsystem Setup
After=multi-user.target

[Service]
Type=oneshot
WorkingDirectory=/home/root
ExecStart=/usr/bin/audiosetup.sh > /dev/null

[Install]
WantedBy=multi-user.target

