# layers/meta-middleby-distro/recipes-core/udev/udev-extraconf/media-automount@.service
[Unit]
Description=Mount media on %i

[Service]
Type=forking
GuessMainPID=no
ExecStart=/etc/udev/scripts/media-automount.sh %I
